# clawdocs.service
[Unit]
Description=ClawDocs - OpenClaw Agent Documentation System
After=network.target

[Service]
Type=simple
User=setrox
Group=setrox
WorkingDirectory=/home/setrox/clawdocs
ExecStart=/usr/bin/npm run serve
Restart=on-failure
RestartSec=5
Environment="NODE_ENV=production"
Environment="PORT=4504"
Environment="FRONTEND_PORT=3504"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/setrox/clawdocs

[Install]
WantedBy=multi-user.target
